# deploy/systemd/bzr-backend.service
[Unit]
Description=BZR Backend API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/bzr-backend
EnvironmentFile=-/var/www/bzr-backend/.env
Environment=NODE_ENV=production
Environment=INGESTER_SUPERVISOR=systemd
ExecStart=/var/www/bzr-backend/scripts/start-backend.sh
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
